# scripts/crontab
SHELL=/bin/bash

*/5 * * * * /bin/sync-conda-envs.sh
